# docker/service/mongod/run
#!/bin/bash

chown -R mongodb /data/configdb /data/db
exec setuidgid mongodb mongod
