# docker/service/prairielearn/run
#!/bin/bash

cd /pl/backend
exec setuidgid pl node server.js /config.json
