# docker/service/prairielearn/run
#!/bin/bash

cd /pl
until psql -U postgres -c "select 1" -d postgres; do
    sleep 1
done
exec setuidgid pl node server.js /config.json
