# docker/service/postgresql/run
#!/bin/bash

[[ -f /.pgvars ]] && source /.pgvars

exec /docker-entrypoint.sh postgres
